# shellscript.service
[Unit]
Description=My Shell Script

[Service]
ExecStart=/usr/bin/script.sh
ExecStop=/usr/bin/shutdownscript.sh

[Install]
WantedBy=multi-user.target
